# systemd/flametech-bridge.service
[Unit]
Description=FlameTech iFlame MQTT Bridge for Home Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bmacdonald3
WorkingDirectory=/home/bmacdonald3
ExecStart=/usr/bin/python3 /home/bmacdonald3/flametech-ha-bridge/src/flametech_mqtt_bridge.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
